# custom/etc/skel/.config/autostart/solaar.desktop
[Desktop Entry]
Type=Application
Name=Solaar
Exec=/usr/bin/solaar --window=hide
Hidden=false
NoDisplay=true
X-GNOME-Autostart-enabled=true
