# system_files/etc/skel/.config/autostart/app.drey.Damask.desktop
[Desktop Entry]
Type=Application
Name=app.drey.Damask
Exec=flatpak run --command=damask app.drey.Damask --background
X-Flatpak=app.drey.Damask
