# custom/1password/systemd/onepassword-cli-perms.path
[Unit]
Description=Watch Homebrew Caskroom for 1Password CLI (re)installs

[Path]
PathModified=/home/linuxbrew/.linuxbrew/Caskroom/1password-cli
Unit=onepassword-cli-perms.service

[Install]
WantedBy=multi-user.target
