# custom/1password/systemd/onepassword-cli-perms.service
[Unit]
Description=Fix 1Password CLI (op) ownership and setgid for desktop app integration
Documentation=https://developer.1password.com/docs/cli/app-integration/

[Service]
Type=oneshot
ExecStart=/usr/libexec/onepassword-cli-perms
